# build/overlay/etc/supervisor.d/30-wireplumber.ini
[program:wireplumber]
priority=30
# logging
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
# options
command=/usr/bin/wireplumber
stopasgroup=true
stopsignal=INT
# start
autostart=false
autorestart=true
